# scripts/com.pictova.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pictova.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yoldaolmak/Projects/Pictova/.venv/bin/python3</string>
        <string>-m</string>
        <string>pictova.app.cli</string>
        <string>serve</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8040</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yoldaolmak/Projects/Pictova/src</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/yoldaolmak/Projects/Pictova/src</string>
        <key>YO_VISUAL_MEMORY_DB</key>
        <string>/Users/yoldaolmak/Downloads/YO_OS_VIL/data/visual_memory.db</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/pictova-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/pictova-server-err.log</string>
</dict>
</plist>
